# src/zabbix_web/build/modules/build/files/etc/supervisor/conf.d/zabbix_web.conf
[program:zabbix_server]
command = /usr/local/bin/dfg /var/run/zabbix/zabbix_server.pid /usr/sbin/zabbix_server -c /etc/zabbix/zabbix_server.conf
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
