# src/zabbix_web/build/modules/build/files/etc/supervisor/conf.d/zabbix_web_stdout.conf
[program:zabbix_serverstdout]
command = /usr/bin/tail -F /var/log/zabbix/stdout.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
